# user/ports/etc/init.d/bash
#!/sbin/openrc-run

name="bash"
command="/bin/bash"
command_args="$BASH_OPTS"
pidfile="/run/bash.pid"
