# user/base/etc/sv/aether-xfce/run
#!/bin/sh

exec 2>&1

service_dir=/run/runit/runsvdir/current
[ -d "$service_dir" ] || service_dir=/etc/runit/runsvdir/current

for service in dbus elogind polkitd seatd udevd; do
    [ -d "$service_dir/$service" ] || exit 1
    sv check "$service_dir/$service" >/dev/null || exit 1
done

for service in agetty-tty1 agetty-tty2 agetty-tty3 agetty-tty4 agetty-tty5 agetty-tty6; do
    [ -d "$service_dir/$service" ] || continue
    sv down "$service_dir/$service" >/dev/null 2>&1 || true
done

udevadm settle || exit 1

exec runuser -l root -c /usr/lib/aether/start-xfce4-session
